# dist/init/groupthink.conf
# groupthink init script
# /etc/init/groupthink.conf

description "groupthink init"
start on runlevel [2345]
stop on runlevel [!2345]

chdir /home/ubuntu/go/src/github.com/aodin/groupthink/
exec ./groupthink -c ./settings.json -l=./logs/requests.log
